# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log

[program:nextjs]
command=node server.js
directory=/app
user=nextjs
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
stderr_logfile=/var/log/nextjs/error.log
stdout_logfile=/var/log/nextjs/access.log

